# supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:sshd]  
command=/usr/sbin/sshd -D

[program:Xvfb]
command=/usr/bin/Xvfb

[program:node]
command=/usr/bin/node /src/server.js 
