# supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:sshd]  
command=/usr/sbin/sshd -D
command=/usr/bin/Xvfb
command=/usr/npm
